# eyejo/supervisord.conf
[program:uwsgi]
user=root
command=uwsgi --ini /root/eyejo/uwsgi.ini
autostart=true
autorestart=true
redirect_stderr=true
stopasgroup=true
killasgroup=true
stopsignal=QUIT

[supervisord]
